# systemd/cortex-core.service
[Unit]
Description=Cortex Core — Wearable AI Memory
After=bluetooth.target sound.target network-online.target
Wants=bluetooth.target sound.target network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/home/turfptax/cortex-core/src
ExecStart=/usr/bin/python3 /home/turfptax/cortex-core/src/main.py
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
